# deploy/maxadmin.service
[Unit]
Description=Velvet Mebel — веб-админка MAX-бота (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/VelvetMebel
ExecStart=/root/VelvetMebel/venv/bin/python -m admin.main
Restart=always
RestartSec=10
TimeoutStopSec=15
EnvironmentFile=/root/VelvetMebel/.env

[Install]
WantedBy=multi-user.target
